# crontab
#run python script every minutes
#https://stackoverflow.com/a/61513773/9599554
* * * * *  /usr/bin/python3 /app/main.py > /proc/1/fd/1 2>/proc/1/fd/2
